# init/systemd/hooky.service
[Unit]
Description=Hooky webhook server
Documentation=https://github.com/virtuallytd/hooky
After=network.target

[Service]
Type=simple
User=hooky
Group=hooky
ExecStart=/usr/local/bin/hooky -hooks /etc/hooky/hooks.yaml -log-format json
Restart=on-failure
RestartSec=5s

# Load secrets and environment variables from a file.
# Create this from .env.example and edit with your real values.
EnvironmentFile=/etc/hooky/.env

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# read-only access to all home directories; /home/hooky is overridden to writable
# below so docker login can store credentials in ~/.docker/config.json.
# Create this directory if it doesn't exist: sudo mkdir -p /home/hooky && sudo chown hooky:hooky /home/hooky
ProtectHome=read-only
ReadWritePaths=/etc/hooky /var/log/hooky /opt/hooky /home/hooky

[Install]
WantedBy=multi-user.target
